# services/blurbcloud/blurbcloud.service
[Unit]
Description=Serve Blurbs
After=litestream-restore.service
Requires=litestream-restore.service

[Service]
ExecStart=/home/blurbcloud/blurbcloud/blurbcloud
WorkingDirectory=/home/blurbcloud/blurbcloud

[Install]
WantedBy=blurbcloud.target
